# stata-vnc/novnc.conf
[program:novnc]
command=/usr/share/novnc/utils/launch.sh --listen 6901 --vnc 127.0.0.1:5901 --web /usr/share/novnc
