# stata-vnc/vncserver.conf
[program:vncserver]
command=/usr/bin/startvnc.sh
